# etc/init/tund.conf
description "tund"

start on net-device-up
stop on runlevel [016]

respawn

exec /usr/local/bin/tund
